# roles/ansible-role-openhab/templates/socat.service.j2
[Unit]
After=network.service

[Service]
Type=simple
PIDFile=/var/run/socat.pid
ExecStart=/usr/local/bin/socat.sh
User=root
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
