# roles/ansible-role-gpsd/templates/gpsd.service.j2
[Unit]
Description=GPS Service
Requires=gpsd.socket

[Service]
EnvironmentFile={{ gpsd_config }}
ExecStart=/usr/sbin/gpsd -F /run/gpsd.sock -N $GPSD_OPTIONS $DEVICES
